# minuspod/optional-transcriber/minuspod-transcriber.container
[Unit]
Description=OpenVINO Model Server (Whisper API) – Optional
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Requires=minuspod.service
After=minuspod.service

[Container]
Image=docker.io/openvino/model_server:latest-gpu
ContainerName=minuspod-transcriber
HostName=minuspod-transcriber
Network=container:minuspod
User=0:0
Annotation=run.oci.keep_original_groups=1
AddDevice=/dev/dri:/dev/dri
Volume=%h/openvino/models:/models:rw
Exec=--config_path /models/config.json --rest_port 8001
Environment=TZ=America/New_York
AutoUpdate=registry
NoNewPrivileges=yes
CPUQuota=300%
MemoryMax=4G
MemoryReservation=2G

[Service]
Restart=on-failure
ExecStartPre=/usr/bin/sleep 5

[Install]
WantedBy=default.target
